# deployment/scoopcam.service
[Unit]
Description=scoopcam.js - making your environment variables rad
Documentation=https://example.com
After=network-online.target

[Service]
Environment=
Type=simple
User=pi
WorkingDirectory=/home/pi/picamera/CameraApp
ExecStartPre=/usr/bin/git pull origin master
ExecStartPre=/home/pi/.nvm/versions/node/v10.6.0/bin/npm install
ExecStart=/home/pi/.nvm/versions/node/v10.6.0/bin/node /home/pi/picamera/CameraApp/app.js
Restart=on-failure

[Install]
WantedBy=multi-user.target